# certbot/files/certbot.service
[Unit]
Description=Lets Encrypt Renewal

[Service]
Type=oneshot
ExecStart=/root/certbot/get-cert.sh